# supervisord.conf
[program:node]
command = node app.js -d ./dict-revised.unicode.sqlite3
directory = /home/dotcloud/current
