# supervisord.conf
[program:node]
command = node app.js -d ./development.unicode.sqlite3
directory = /home/dotcloud/current
